# app/crontab
* * * * * echo "$(date)" >> /tmp/cron.log
